# scripts/podman/openclaw.container.in
# QCortex gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into qcortex's ~/.config/containers/systemd/
# {{QCORTEX_HOME}} is replaced at install time.

[Unit]
Description=QCortex gateway (rootless Podman)

[Container]
Image=qcortex:local
ContainerName=qcortex
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{QCORTEX_HOME}}/.qcortex:/home/node/.qcortex
EnvironmentFile={{QCORTEX_HOME}}/.qcortex/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
